# tradovate-dashboard.service
[Unit]
Description=Tradovate Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/Tradovate-Bot
ExecStart=/usr/bin/python3 /home/user/Tradovate-Bot/dashboard.py
Restart=always
RestartSec=10
StandardOutput=append:/home/user/Tradovate-Bot/dashboard.log
StandardError=append:/home/user/Tradovate-Bot/dashboard.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
