# tradovate-dashboard.service
[Unit]
Description=Tradovate Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/tradovate-bot
ExecStart=/root/tradovate-bot/venv/bin/python /root/tradovate-bot/dashboard.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradovate-dashboard
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
