# tradovate-agent.service
[Unit]
Description=Tradovate Bot VPS Agent (GitHub-based remote control)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Tradovate-Bot
ExecStart=/usr/bin/python3 /root/Tradovate-Bot/vps_agent.py
Restart=always
RestartSec=10
StandardOutput=append:/root/Tradovate-Bot/agent.log
StandardError=append:/root/Tradovate-Bot/agent.log
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
