# tradovate-bot.service
[Unit]
Description=Tradovate Trading Bot
After=network-online.target
Wants=network-online.target
StartLimitBurst=10
StartLimitIntervalSec=600

[Service]
Type=simple
WorkingDirectory=/root/tradovate-bot
ExecStart=/root/tradovate-bot/venv/bin/python bot.py --live
EnvironmentFile=/root/tradovate-bot/.env
ExecStopPost=/bin/bash /root/tradovate-bot/alert.sh tradovate-bot

# Always restart — the bot has its own daily sleep loop.
# If it exits unexpectedly, systemd brings it back within 30s.
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tradovate-bot
Environment=PYTHONUNBUFFERED=1

# Security
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
